# supervisor_conf/nginx.conf
[program:php-fpm8.3]
command=/usr/sbin/php-fpm8.3  --nodaemonize --fpm-config /etc/php/8.3/fpm/php-fpm.conf
autorestart=true
